# docker/service/app/run
#!/bin/bash

pushd /app
FLASK_APP=app.py FLASK_DEBUG=True exec flask run --host=0.0.0.0
popd
